# datacollector.service.d/datacollector.conf
[Service]
Environment="PYTHONOPTIMIZE=2"
